# ansible/roles/gitea/templates/etc/services.d/gitea.ini.j2
[program:gitea]
directory=%(ENV_HOME)s/gitea
command=%(ENV_HOME)s/gitea/gitea web
startsecs=30
autorestart=yes
